# conf/cron/crontab
*/600 * * * * /usr/local/bin/python3.2 ~/webapps/viddle/crawler/miner.py
